# services/s6-rc.d/apache/finish
#!/bin/sh
# Apache finish script
# We use /bin/sh to ensure simple logging to stdout
echo "[apache] finish: stopped with exit code $1 and signal $2"
