# services/apache/run
#!/bin/sh
# Load environment variables from s6-overlay
# s6-overlay v3 stores them in /var/run/s6/container_environment/
if [ -d /var/run/s6/container_environment ]; then
    for file in /var/run/s6/container_environment/*; do
        if [ -f "$file" ] && [ -r "$file" ]; then
            export "$(basename "$file")=$(cat "$file")"
        fi
    done
fi
# Legacy location for older s6-overlay versions
if [ -d /container/environment ]; then
    for file in /container/environment/*; do
        if [ -f "$file" ] && [ -r "$file" ]; then
            export "$(basename "$file")=$(cat "$file")"
        fi
    done
fi
exec apache2-foreground

